# docker/demo/config/supervisor/server-3.conf
[program:car]
command=tail -f /dev/null

[program:boat]
command=tail -f /dev/null

[program:plane]
command=tail -f /dev/null

[program:train]
command=tail -f /dev/null
