# docker/demo/config/supervisor/server-1.conf
[program:pony]
command=tail -f /dev/null

[program:unicorn]
command=tail -f /dev/null

[program:axolotl]
command=tail -f /dev/null

[program:elephpant]
command=tail -f /dev/null
